# myco.service
[Unit]
Description=Mycelium Web Terminal
After=network.target

[Service]
Type=simple
User=kkrazy
Group=kkrazy
Environment=TLS_CERT_PATH=/home/kkrazy/myco/.tls/cert.pem
Environment=TLS_KEY_PATH=/home/kkrazy/myco/.tls/key.pem
EnvironmentFile=/home/kkrazy/myco/.env
WorkingDirectory=/home/kkrazy/myco
ExecStart=/home/kkrazy/.local/node/bin/node server/src/index.js
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
